# Rasp/autoboot/gilberto.service
[Unit]
After=idle

[Service]
ExecStart=/home/pi/AVAST/Rasp/autoboot/start-gilberto.sh

[Install]
WantedBy=default.target
